# bandwidth_occupier.service
[Unit]
Description=Oracle keep-alive bandwidth occupier
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/oalive/oalive.conf
ExecStart=/bin/sh /usr/local/bin/bandwidth_occupier.sh
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
