# cpu-limit.service
[Unit]
Description=Oracle keep-alive CPU occupier
After=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/oalive/oalive.conf
ExecStart=/bin/sh /usr/local/bin/cpu-limit.sh
Restart=always
RestartSec=30s
KillMode=control-group

[Install]
WantedBy=multi-user.target
